# deployment/ubuntu24/mr-moneybags.service
[Unit]
Description=Mr MoneyBags API Server
After=network.target postgresql.service

[Service]
User=mrmb
Group=mrmb
WorkingDirectory=/opt/mr-moneybags
EnvironmentFile=/opt/mr-moneybags/.env
Environment=NODE_ENV=production
ExecStart=/usr/bin/node server-modular.js
Restart=on-failure
RestartSec=5

# Security hardening
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=full
ProtectHome=true

[Install]
WantedBy=multi-user.target
